# includes.container/etc/xdg/autostart/clipboard-sync.desktop
[Desktop Entry]
Name=Wayland/X Clipboard Sync
Exec=systemctl --user start clipboard-sync
Type=Application
